# etc/local.d/beep.start
echo -e '\a'

gateway=$(/bin/ip r | awk '/default/ { print $3 }' 2>/dev/null)
test "$gateway" && ping -qc 1 $gateway >/dev/null && echo -e '\a'
