# etc/local.d/startup-beep.start
echo -e '\a'

ping -qc $(ip r | awk '/default/ { print $3 }') >/dev/null && echo -e '\a'
